# systemd/irccd.service.in
[Unit]
Description=IRC Client Daemon
After=network.target

[Service]
ExecStart=@PATH@
User=@UID@
Group=@GID@

[Install]
WantedBy=multi-user.target
