# systemd/irccd.service
[Unit]
Description=IRC Client Daemon
After=network.target

[Service]
ExecStart=@PATH@

[Install]
WantedBy=multi-user.target
