# systemd/irccd.service.in
[Unit]
Description=IRC Client Daemon
After=network.target

[Service]
ExecStart=@PATH@
User=@IRCCD_WITH_UID@
Group=@IRCCD_WITH_GID@

[Install]
WantedBy=multi-user.target
